# my-crawlers.service
[Unit]
Description=My Crawlers Tracker Service
After=network.target

[Service]
Type=simple
User=azureuser
WorkingDirectory=/home/azureuser/my-crawlers
ExecStart=/usr/bin/node tracker/index.js
Restart=always
Environment=NODE_ENV=production
EnvironmentFile=/home/azureuser/my-crawlers/.env

[Install]
WantedBy=multi-user.target
